# rootfs/usr/lib/systemd/system/pios-swap.service
[Unit]
Description=Pi AB Swapper
After=sysinit.target network-online.target
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=/usr/lib/pi-updater/pios_swap.sh
