# rootfs/var/lib/dpkg/info/pi-bluetooth.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bthelper@.service
/lib/systemd/system/hciuart.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/90-pi-bluetooth.rules
/usr
/usr/bin
/usr/bin/bthelper
/usr/bin/btuart
/usr/share
/usr/share/doc
/usr/share/doc/pi-bluetooth
/usr/share/doc/pi-bluetooth/changelog.gz
/usr/share/doc/pi-bluetooth/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/pi-bluetooth
