# rootfs/var/lib/dpkg/info/raspi-config.conffiles
/etc/default/cpu_governor
/etc/init.d/raspi-config
/etc/sudoers.d/010_proxy
