# rootfs/var/lib/dpkg/info/raspberrypi-net-mods.list
/.
/etc
/etc/profile.d
/etc/profile.d/wifi-check.sh
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/raspberrypi-net-mods.service
/usr
/usr/lib
/usr/lib/raspberrypi-net-mods
/usr/lib/raspberrypi-net-mods/wpa_copy
/usr/share
/usr/share/doc
/usr/share/doc/raspberrypi-net-mods
/usr/share/doc/raspberrypi-net-mods/changelog.gz
