# rootfs/lib/systemd/system/NetworkManager.service.d/10-dhcpcd.conf
[Unit]
After=dhcpcd.service
Conflicts=dhcpcd.service
