# rootfs/usr/lib/systemd/system/multi-user.target.wants/image-importer.service
[Unit]
Description=Image Importer
After=multi-user.target

[Service]
Type=simple
ExecStart=/scripts/image-importer.sh
