# rootfs/usr/lib/systemd/system/multi-user.target.wants/ntp-updater.service
[Unit]
Description=NTP Updater
After=sysinit.target network-online.target
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=/scripts/ntp-updater.sh